# roles/prometheus/templates/node-exporter.service
[Unit]
Description = Node Exporter

[Service]

ExecStart=/bin/node /app/index.js
SyslogIdentifier=backend

[Install]
WantedBy=multi-user.target